# example-mqtt/activemq.conf
[program:activemq]
;user=nobody
command=/opt/activemq/bin/activemq console
redirect_stderr=true
stdout_logfile=/var/log/activemq.log
startretries=5
